# deb/vc_wifitracker/lib/systemd/system/vc-wifitracking.service
[Unit]
Description=VC WiFi Tracking Service
After=network.target

[Service]
Type=simple
User=vc-wifitracking
Group=vc-wifitracking
WorkingDirectory=/opt/vc-wifitracking
ExecStartPre=/usr/sbin/ip link set wlan0 down
ExecStartPre=/usr/sbin/iw dev wlan0 set type monitor
ExecStartPre=/usr/sbin/ip link set wlan0 up
ExecStart=/usr/bin/python3 /opt/vc-wifitracking/app.py
ExecStopPost=/usr/sbin/ip link set wlan0 down
ExecStopPost=/usr/sbin/iw dev wlan0 set type managed
ExecStopPost=/usr/sbin/ip link set wlan0 up

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

Restart=on-failure

[Install]
WantedBy=multi-user.target
